# media-tv/plex-media-server/files/plex-media-server
#!/sbin/runscript

# Copyright 2015 Jonathan Vasquez <jvasquez1011@gmail.com>
# Distributed under the terms of the GNU General Public License v2

PMS_DISPLAY_NAME="Plex Media Server"
PMS_DASH_NAME="${PMS_DASH_NAME:-plex-media-server}"
PMS_SCRIPT_NAME="start_pms"
PMS_PATH="${PMS_PATH:-/usr/sbin}"
PMS_BINARY="${PMS_BINARY:-$PMS_PATH/$PMS_SCRIPT_NAME}"
PMS_PIDFILE="${PMS_PIDFILE:-/var/run/$PMS_DASH_NAME.pid}"

depend() {
        need avahi-daemon
}

start() {
	ebegin "Starting ${PMS_DISPLAY_NAME}"
	start-stop-daemon -m --pidfile "${PMS_PIDFILE}" -b -u plex -N -5 --start --exec "${PMS_BINARY}"
	eend $?
}

stop() {
	ebegin "Stopping ${PMS_DISPLAY_NAME}"
	start-stop-daemon --stop --exec "${PMS_BINARY}" --pidfile "${PMS_PIDFILE}"
	eend $?
}

reload() {
	ebegin "Reloading ${PMS_DISPLAY_NAME}"
	start-stop-daemon --signal HUP --exec "${PMS_BINARY}" --pidfile "${PMS_PIDFILE}"
	eend $?
}
