# app-admin/lsyncd/files/init.d/lsyncd
#!/sbin/runscript

# Copyright 2015 Jonathan Vasquez <jvasquez1011@gmail.com>
# Distributed under the terms of the GNU General Public License v2

LSYNCD_NAME="lsyncd"
LSYNCD_PATH="/usr/bin"
LSYNCD_BINARY="${LSYNCD_PATH}/${LSYNCD_NAME}"

start() {
	ebegin "Starting ${LSYNCD_NAME}"

	start-stop-daemon --start --exec "${LSYNCD_BINARY}" \
	--make-pidfile --pidfile "${LSYNCD_PIDFILE}" \
	--background \
	--stdout "${LSYNCD_STDLOG}" \
	--stderr "${LSYNCD_ERRLOG}" \
	-- ${LSYNCD_OPTS}

	eend $?
}

stop() {
	ebegin "Stopping ${LSYNCD_NAME}"

	start-stop-daemon --stop --exec "${LSYNCD_BINARY}" --pidfile "${LSYNCD_PIDFILE}"

	eend $?
}
